# packages/flem-tool/flem_tool-0.2.0.tar.gz/flem_tool-0.2.0/flem.service
[Unit]
Description=FLEM Tool Service


[Service]
ExecStart=flem
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
Type=notify


[Install]
WantedBy=default.target